# Ansible/crontab
    */5 * * * * root sh /home/vagrant/jenkins-job.sh